# server/usr/lib/systemd/system/pulp_workers.service
[Unit]
Description=Pulp Celery Workers
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/var/run/pulp/
ExecStart=/usr/bin/python -m  pulp.server.async.manage_workers start
ExecStop=/usr/bin/python -m  pulp.server.async.manage_workers stop

[Install]
WantedBy=multi-user.target
